# ansible/roles/supervisor/files/conf.d/s136.okserver.org/rt.conf
[program:rt]
command=/opt/rt4/sbin/rt-server.fcgi --port 9000
autostart=true
autorestart=true
user=www-data
group=www-data
stdout_logfile=/var/log/rt.log
redirect_stderr=true
stopsignal=QUIT
