# ansible/roles/scripts/files/opt/scripts/mail/cron.d/mailman_weekly_stats
0 0 * * 0 root /opt/scripts/mailman-weekly-stats.py
