# etc/supervisor/conf.d/malaria.bibsoup.net.conf
[program:malaria.bibsoup.net]
command=/home/okfn/var/srvc/%(program_name)s/bin/gunicorn -w 4 -b 127.0.0.1:5052 -n %(program_name)s bibserver.web:app
environment=PATH='/home/okfn/var/srvc/malaria.bibsoup.net/bin:/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin'
user=okfn
directory=/home/okfn/var/srvc/%(program_name)s/src/bibserver
stdout_logfile=/var/log/supervisor/%(program_name)s-access.log
stderr_logfile=/var/log/supervisor/%(program_name)s-error.log
autostart=true
