# roles/nagios-server/files/etc/cron.d/pagerduty_nagios
* * * * * root /usr/local/bin/pagerduty_nagios.pl flush
