# roles/scripts/files/opt/scripts/graphite_collectors/cron.d/exim_stats_graphite
*/5 * * * * root /opt/scripts/graphite_collectors/exim_stats.py
