# roles/scripts/files/opt/scripts/graphite_collectors/cron.d/mailman_stats_graphite
*/5 * * * * root /opt/graphite_scripts/mailman_stats.py
