# roles/graphite_collectors/files/cron.d/exim_stats_graphite.cron
*/5 * * * * root /opt/graphite_scripts/exim_stats.py
