# ansible/roles/scripts/files/opt/scripts/graphite_collectors/cron.d/nginx_stats_graphite
*/5 * * * * root /opt/scripts/graphite_collectors/nginx_stats.py
