# ansible/roles/scripts/files/opt/scripts/graphite_collectors/cron.d/postfix_spam_stats_graphite
*/5 * * * * root /opt/scripts/graphite_collectors/postfix_spam_stats.py
