# ansible/roles/scripts/files/opt/scripts/security-updates/cron.d/security-updates
05 */2 * * * root /opt/scripts/security-updates
