# ansible/roles/scripts/files/opt/scripts/backup/cron.d/psql_backup
30 02 * * * root /opt/scripts/backup/psql_backup.sh &> /var/log/psql_backup.log && curl https://nosnch.in/304f7bf902 &> /dev/null
