# roles/scripts/files/opt/scripts/cron.d/mail_forwards_merge
*/15 * * * * root /usr/bin/python2.7 /opt/scripts/mail_forwards_merge.py
