# ansible/roles/db-backup/files/opt/scripts/backup/cron.d/mysql_backup
30 01 * * * root /opt/scripts/backup/mysql_backup.sh &> /var/log/mysql_backup.log && curl https://nosnch.in/ecc09ca67b &> /dev/null
