# ansible/roles/db-backup/files/etc/cron.d/psql_backup
30 02 * * * root /opt/scripts/backup/psql_backup.sh &> /var/log/psql_backup.log 
